# build/ccm-agent.service
[Unit]
Description=The Choria Configuration Manager Agent
After=network.target

[Service]
User=root
Group=root
ExecStart=/usr/bin/ccm agent --config /etc/choria/ccm/agent.yaml
KillMode=process
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target